# aports/tautulli/tautulli.initd
#!/sbin/openrc-run

# fallback values for /etc/conf.d/tautulli
: ${command_user:=tautulli:tautulli}
: ${directory:=/var/lib/tautulli}

name=tautulli
description="Tautulli daemon"
command="/usr/bin/python3 /usr/lib/tautulli/Tautulli.py"
command_args="--nofork --nolaunch --datadir $directory"

depend() {
	need net
	use dns logger
}

start_pre() {
	checkpath -d -o "$command_user" "$directory"
}
