# aports/sabnzbd/sabnzbd.initd
#!/sbin/openrc-run

# fallback values for /etc/conf.d/sabnzbd
: ${command_user:=sabnzbd:sabnzbd}
: ${directory:=/var/lib/sabnzbd}
: ${umask:=0002}

name=sabnzbd
description="SABnzbd daemon"
command="/usr/bin/python3"
command_args="-m SABnzbd --config-file $directory/sabnzbd.ini $command_args"

depend() {
	need net
	use dns logger
}

start_pre() {
	checkpath -d -o "$command_user" "$directory"
}
